# scripts/com.lampese.codex-switcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Unique identifier — must match filename -->
  <key>Label</key>
  <string>com.lampese.codex-switcher</string>

  <!-- Path to the app binary inside the .app bundle -->
  <key>ProgramArguments</key>
  <array>
    <string>APP_PATH</string>
  </array>

  <!-- Start at login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Keep alive: restart if it crashes -->
  <key>KeepAlive</key>
  <false/>

  <!-- Throttle restarts to prevent tight crash loops -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Silence stdout/stderr (it's a GUI app — no terminal output needed) -->
  <key>StandardOutPath</key>
  <string>/tmp/codex-switcher.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/codex-switcher-error.log</string>
</dict>
</plist>
